# src/launchd/ai.brain.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.brain.indexer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Shared/ai-brain/venv/bin/python</string>
        <string>/Users/Shared/ai-brain/scripts/indexer.py</string>
        <string>--config</string>
        <string>/Users/Shared/ai-brain/indexer_config.yaml</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Shared/ai-brain/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/Shared/ai-brain/logs/indexer.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Shared/ai-brain/logs/indexer.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/Shared/ai-brain</string>
</dict>
</plist>
